# logsentinel.service
# LogSentinel — systemd unit for the supervisor-managed shard fleet.
#
# Install:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin logsentinel
#   sudo mkdir -p /opt/logsentinel /var/lib/logsentinel
#   sudo cp logsentinel*.py rules.yaml /opt/logsentinel/
#   sudo chown -R logsentinel:logsentinel /var/lib/logsentinel
#   sudo cp logsentinel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now logsentinel
#
# Inspect:
#   systemctl status logsentinel
#   journalctl -u logsentinel -f
#   curl localhost:9200/metrics        # shard 0
#   curl localhost:9204/metrics        # aggregator (base + shards, here 4)
#
# Tune the shard count, ports, log paths, and dedup window in ExecStart below.

[Unit]
Description=LogSentinel security log analysis (sharded streaming fleet)
Documentation=file:/opt/logsentinel/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=logsentinel
Group=logsentinel

# The supervisor launches N shards + the aggregator and keeps them alive.
# It writes all per-shard outputs, checkpoints, and the merged stream under
# --work-dir. Adjust the log paths and --shards to your host.
ExecStart=/usr/bin/python3 /opt/logsentinel/logsentinel_supervisor.py \
    /var/log/auth.log /var/log/nginx/access.log \
    --shards 4 \
    --base-metrics-port 9200 \
    --work-dir /var/lib/logsentinel \
    --merged-output /var/lib/logsentinel/merged.jsonl \
    --dedup-window 60 \
    --rules /opt/logsentinel/rules.yaml

# The supervisor traps SIGTERM and shuts its children down cleanly; give it
# time to do so before systemd force-kills the group.
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resilience: restart the whole supervisor if it ever exits unexpectedly.
# (The supervisor already restarts individual shards; this covers the
# supervisor process itself.)
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# --- Security hardening ---
# LogSentinel only needs to READ log files and WRITE to its state dir.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
# Only these paths are writable; everything else on / is read-only to the unit.
ReadWritePaths=/var/lib/logsentinel
# Log sources are read-only; add any others your ExecStart follows.
ReadOnlyPaths=/var/log
# Restrict network to loopback (metrics endpoints) + normal egress if a future
# alert sink needs it. Tighten to AF_UNIX/AF_INET if you don't need IPv6.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource guardrails (a sharded fleet is CPU-bound; memory stays bounded).
# Tune to your host; these are conservative starting points for 4 shards.
MemoryMax=1G
TasksMax=64

[Install]
WantedBy=multi-user.target
